# crontab
*/2 * * * * root cd /app && python3 scrapy.py >> /var/log/cron.log 2>&1
